# deploy/blackjack.service
[Unit]
Description=Blackjack Simulator Web App
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/blackjack
Environment=PATH=/opt/blackjack/venv/bin:/usr/bin:/bin
Environment=PORT=8000
Environment=CORS_ORIGINS=*
ExecStart=/opt/blackjack/venv/bin/uvicorn web.api:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
